# deploy/com.david.gloss.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Gloss server LaunchAgent for the Mac Studio (run as user david).

  Install:
    1. Replace __GLOSS_AUTH_TOKEN__ below (generate: openssl rand -hex 32)
    2. cp deploy/com.david.gloss.plist ~/Library/LaunchAgents/
    3. launchctl bootstrap gui/501 ~/Library/LaunchAgents/com.david.gloss.plist

  See DEPLOY.md for the full runbook.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.david.gloss</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/david/.bun/bin/bun</string>
    <string>/Users/david/gloss/src/cli.ts</string>
    <string>serve</string>
    <string>--port</string>
    <string>3456</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/david/gloss</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>GLOSS_REMOTE</key>
    <string>1</string>
    <key>GLOSS_BIND_HOST</key>
    <string>100.109.110.36</string>
    <key>GLOSS_AUTH_TOKEN</key>
    <string>__GLOSS_AUTH_TOKEN__</string>
    <key>GLOSS_DISABLE_OS_ENDPOINTS</key>
    <string>1</string>
    <!-- Studio-native logs + synced laptop logs, attributed by root -->
    <key>GLOSS_PROJECTS_ROOTS</key>
    <string>studio=/Users/david/.claude/projects,mbp=/Users/david/.claude/projects-mbp</string>
    <key>HOME</key>
    <string>/Users/david</string>
    <!-- bun's transpiler cache (~/Library/Caches/bun) has hung module loads
         on this machine; the cache only speeds up cold starts, so disable -->
    <key>BUN_RUNTIME_TRANSPILER_CACHE_PATH</key>
    <string>0</string>
    <key>PATH</key>
    <string>/Users/david/.bun/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>/Users/david/Library/Logs/gloss.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/david/Library/Logs/gloss.err.log</string>
</dict>
</plist>
